# systems/omega/etc/cron.d/charybdis
*/1 * * * * charybdis /usr/local/bin/ircd.cron
