# systems/illium/etc/cron.d/forum
*/1 * * * * www-data php /var/www/FUDforum/scripts/cron.php >/dev/null 2>&1
